# inittab
# (c) Copyright 2019-2020, James Stevens ... see LICENSE for details
# Alternative license arrangements possible, contact me for more information

::respawn:/sbin/syslogd -n -O -
::respawn:/opt/start_epprest 1
::respawn:/opt/start_epprest 2
::respawn:/opt/start_nginx
